# deploy/vizthinker-backend.service
[Unit]
Description=VizThinker Backend Service
After=network.target

[Service]
Type=simple
User=jemmy
Group=jemmy
WorkingDirectory=/home/jemmy/vizthinker
Environment=PATH=/home/jemmy/vizthinker/.venv/bin
ExecStart=/home/jemmy/vizthinker/.venv/bin/uvicorn server.main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 